# timers/_each-1mn-daily.timer
[Unit]
Description=Timer each minutes daily

[Timer]
AccuracySec=1min
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
