# timers/_twice-daily.timer
[Unit]
Description=Timer twice daily

[Timer]
AccuracySec=1min
OnCalendar=*-*-* 00,12:00:00

[Install]
WantedBy=timers.target
